# infrastructure/c1-link-setup/roles/c1-link-router/files/homelab-router.service
[Unit]
Description=Homelab L3 router glue (lw-c1 direct cable)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/homelab-router-setup.sh

[Install]
WantedBy=multi-user.target
