# monitoring/kiosk-setup/roles/kiosk/templates/kiosk.service.j2
[Unit]
Description=Grafana NOC wall kiosk
After=graphical.target

[Service]
User={{ kiosk_user }}
Environment=DISPLAY=:0
# Needed for a system service to open the user's X session on :0.
Environment=XAUTHORITY=/home/{{ kiosk_user }}/.Xauthority
ExecStart=/usr/local/bin/kiosk-start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target
