# k8s/k3s-setup/files/k3s-backup.timer
[Unit]
Description=Daily k3s state backup at 03:15 UTC

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
