# cloud-storage/systemd/homelab-s3-sync.service
[Unit]
Description=Homelab MinIO → B2 + Tigris + R2 offsite backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kamil-rybacki
Group=kamil-rybacki
# Env file is sourced by the script itself (uses `export VAR=...` shell syntax
# which systemd EnvironmentFile cannot parse). systemd just needs HOME set,
# which it does automatically from User=.
ExecStart=/home/kamil-rybacki/Code/ansible/cloud-storage/06-sync-minio-to-cloud.sh --confirm
StandardOutput=journal
StandardError=journal
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
