# cloud-storage/systemd/homelab-s3-sync.service
[Unit]
Description=Homelab MinIO → B2 + Tigris + R2 offsite backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kamil-rybacki
Group=kamil-rybacki
# %h resolves before User= drops privileges, so it would point at root's home.
# Hardcode absolute path to the actual user's env file.
EnvironmentFile=/home/kamil-rybacki/.config/homelab-s3-remotes.env
ExecStart=/home/kamil-rybacki/Code/ansible/cloud-storage/06-sync-minio-to-cloud.sh --confirm
StandardOutput=journal
StandardError=journal
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
